# deployment/templates/gunicorn.conf.jinja2
; =====================
;  manage.py supervisor
; =====================

[program:gunicorn]
command=/home/ubuntu/.virtualenvs/NCDR-reference/bin/gunicorn ncdr_reference.wsgi -c %(here)s/gunicorn_conf.py
directory=%(here)s/../../
;user=nobody
numprocs=1
stdout_logfile={{ LOG_DIR }}/gunicorn.log
stderr_logfile={{ LOG_DIR }}/gunicorn.log
autostart=true
autorestart=true
startsecs=10
loglevel=DEBUG
